# systemd/contradiction-sweep.service
[Unit]
Description=v0.19 I.3: weekly offline contradiction sweep (canonical vs similar candidates; cheap local discovery)
After=network-online.target

[Service]
Type=oneshot
# C5 (2026-07-03, audit decision 2026-06-14): contradiction judging is JUDGMENT
# work -> Codex via the Windows shim (:18792); local models are never the judge
# (an earlier 3B local judge YES'd 9/9 pairs spuriously, and the v0.27.3
# re-judge measured 78% false positives on local verdicts). This REVERSES the
# v0.27.3 `--judge local` trade: the script preflights the shim and NO-OPs
# (exit 0, outcome=no-op:codex-shim-unreachable) when it is down — at the
# unattended Sun 05:00 slot that usually means the week is SKIPPED, and that is
# the intended trade: a skipped week (false NO) is strictly better than a weak
# local YES polluting the store with spurious advisory stamps. It NEVER falls
# back to local (that would re-introduce the audited misrouting). The
# session-time SessionStart rejudge (storage-cap-check.sh, fires when the shim
# IS up, weekly-throttled) remains the authoritative resolution path, so
# contradiction coverage continues on any week the operator opens a session.
ExecStart=/home/__WSL_USER__/apps/mem0-server/.venv/bin/python %h/apps/mem0-scripts/contradiction-sweep.py --apply --limit 50 --judge codex
StandardOutput=journal
StandardError=journal
