# systemd/contradiction-sweep.timer
[Unit]
Description=Weekly contradiction-sweep timer (Sun 05:00 - after goals-stale-sweep 04:00; LLM-heavy, runs off-hours)

[Timer]
OnCalendar=Sun 05:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
