# systemd/decay-scan.service
[Unit]
Description=Weekly decay + dedup scan (v0.13)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/__WSL_USER__/apps/mem0-server/.venv/bin/python %h/apps/mem0-scripts/decay-scan.py
ExecStartPost=/home/__WSL_USER__/apps/mem0-server/.venv/bin/python %h/apps/mem0-scripts/semantic-dedup.py
