# systemd/goals-stale-sweep.service
[Unit]
Description=v0.17 Phase E: weekly stale-goal sweeper (report-only)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/__WSL_USER__/apps/mem0-server/.venv/bin/python %h/apps/mem0-scripts/goals-stale-sweep.py
StandardOutput=journal
StandardError=journal
