# systemd/mem0.service
[Unit]
Description=mem0 v2.0.4 WSL-native server (v0.12 agentic-memory-stack)
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/apps/mem0-server
# v0.29.3: R4 raw-trace fallback ON by default in production (its episodic
# test-pollution gate was cleared + live-verified in v0.29.2). Surfaces ONE
# semantic-relevant past-episode snippet ONLY when condensed search returns 0
# memories; off-domain prompts abstain (cosine floor). Disable with =0 if needed.
Environment=MEM0_RAW_FALLBACK_ENABLED=1
# v0.29.4 item 1: extend Weibull freshness to the time-sensitive tiers (evidence,
# temporal) on the durable read path — atemporal canonical/stable/insight never decay.
# GENTLE (365d half-life default); demotes only genuinely-stale evidence. Disable with =0.
Environment=MEM0_DURABLE_FRESHNESS_ENABLED=1
# v1.0 R7 (Phase 7A): operator-agnostic default tenant for the /v1/context/bundle
# proactive search. The installer substitutes the __WSL_USER__ sentinel with the
# install user (1-wsl-services.sh systemd sed) so a third-party box queries its
# own tenant; app.py falls back to "default" if unset.
Environment=MEM0_DEFAULT_USER_ID=__WSL_USER__
# MEM-18 (2026-07-03): kill runtime phone-home. mem0 2.0.4 ships PostHog
# telemetry ON by default — mem0/memory/telemetry.py reads MEM0_TELEMETRY and
# treats anything outside ("true","1","yes") as off, so "False" disables it
# (verified against the installed lib source, not the docs). HF_HUB_OFFLINE=1
# is belt-and-braces: the reranker/embedder are HTTP clients to llama-swap (no
# in-process HF load), so this only stops any transitive huggingface-hub
# dependency from ever calling home from the server process — it can't break
# inference because nothing in this process downloads weights.
Environment=MEM0_TELEMETRY=False
Environment=HF_HUB_OFFLINE=1
# v0.19 Phase H key chain: private tmpfs runtime dir + best-effort DPAPI key
# fetch before start. The '-' prefix keeps fresh/plaintext-mode boxes safe:
# the fetch fails soft and canonical_key_provider falls through to plaintext.
RuntimeDirectory=mem0
RuntimeDirectoryMode=0700
ExecStartPre=-%h/apps/mem0-server/dpapi-fetch-key.sh
# Bind address is operator config: loopback by default (safe single-box posture).
# A multi-box deployment where replica clients reach this brain over the network
# (e.g. a tailnet) sets MEM0_BIND=0.0.0.0 in ~/.mem0/stack.env — the deployers
# (deploy.sh / 1-wsl-services.sh) resolve __MEM0_BIND__ from it at install time.
ExecStart=%h/apps/mem0-server/.venv/bin/uvicorn app:app --host __MEM0_BIND__ --port 18791
Restart=on-failure
RestartSec=5
MemoryMax=2G

[Install]
WantedBy=default.target
