# systemd/stack-backup.service
[Unit]
Description=Weekly agentic memory stack snapshot

[Service]
Type=oneshot
ExecStart=/bin/bash %h/apps/mem0-scripts/stack-backup.sh
